# ansible/roles/common/templates/supervisor/reload-dnsmasq.conf
[program:reload-dnsmasq]
user = root
command = /opt/libexec/reload-dnsmasq
redirect_stderr = true
